# dailyPush.cron
0  21  *  *  * /home/elamd/.scripts/dailyPush
